# tradeliy.service.template
[Unit]
Description=Tradeliy Trading Platform
After=network.target mongodb.service postgresql.service

[Service]
Type=simple
User=tradeliy
WorkingDirectory=/home/tradeliy/tradeliy
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=tradeliy
Environment=NODE_ENV=production
Environment=PORT=5000
Environment=PYTHON_API_URL=http://localhost:5001

[Install]
WantedBy=multi-user.target